# crontab
*/1 * * * * py /app/main.py >> /dev/stdout
