# crontab
*/1 * * * * py /app/main.py >> /var/log/cron/app.log 2>&1